# KQL/sigma-generic/Execution/T1053.005/Suspicious_Scheduled_Task_Update.kql
// Title: Suspicious Scheduled Task Update
// Author: Nasreddine Bencherchali (Nextron Systems)
// Date: 2022-12-05
// Level: high
// Description: Detects update to a scheduled task event that contain suspicious keywords.The SecurityEvent table in Microsoft Sentinel contains security-related events from Windows systems. To query the Command field from the embedded XML in the event data, you may need to parse the XML data using Kusto query language functions like parse_xml(). Ensure that you have the necessary data connectors and log sources configured in Sentinel to collect Windows security events.
// Tags: attack.execution, attack.privilege-escalation, attack.persistence, attack.t1053.005
// ================================================================== 

SecurityEvent
| where (TaskContentNew contains "regsvr32" or TaskContentNew contains "rundll32" or TaskContentNew contains "cmd.exe</Command>" or TaskContentNew contains "cmd</Command>" or TaskContentNew contains "<Arguments>/c " or TaskContentNew contains "<Arguments>/k " or TaskContentNew contains "<Arguments>/r " or TaskContentNew contains "powershell" or TaskContentNew contains "pwsh" or TaskContentNew contains "mshta" or TaskContentNew contains "wscript" or TaskContentNew contains "cscript" or TaskContentNew contains "certutil" or TaskContentNew contains "bitsadmin" or TaskContentNew contains "bash.exe" or TaskContentNew contains "bash " or TaskContentNew contains "scrcons" or TaskContentNew contains "wmic " or TaskContentNew contains "wmic.exe" or TaskContentNew contains "forfiles" or TaskContentNew contains "scriptrunner" or TaskContentNew contains "hh.exe") and (TaskContentNew contains "\\AppData\\Local\\Temp\\" or TaskContentNew contains "\\AppData\\Roaming\\" or TaskContentNew contains "\\Users\\Public\\" or TaskContentNew contains "\\WINDOWS\\Temp\\" or TaskContentNew contains "C:\\Temp\\" or TaskContentNew contains "\\Desktop\\" or TaskContentNew contains "\\Downloads\\" or TaskContentNew contains "\\Temporary Internet" or TaskContentNew contains "C:\\ProgramData\\" or TaskContentNew contains "C:\\Perflogs\\")